# cron.d/kg-pbs-logger
# Process yesterdays logs at 4am

0 4 * * *       root    /usr/sbin/kg-pbs-logger -y

