# debian/karaage-cluster-tools.cron.d
# Process yesterdays logs at 4am

0 4 * * *       root    kg-send-usage -y

